# deploy/queue-worker.supervisord.conf
[program:mines-queue-worker]
command=php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --timeout=90
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/mines/queue-worker.log
stopsignal=TERM
